# feeders/xmm7360/cellsignal-xmm7360.timer
[Unit]
Description=Refresh Cell Signal metrics feed

[Timer]
OnBootSec=30
OnUnitActiveSec=2
AccuracySec=1

[Install]
WantedBy=timers.target
