# feeders/xmm7360/cellsignal-xmm7360.service
[Unit]
Description=Publish XMM7360 cellular metrics to /run/cellsignal.json (Cell Signal feeder)
ConditionPathExists=/dev/wwan0xmmrpc0

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cellsignal-feeder-xmm7360 --write
